# server/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn main:app -k uvicorn.workers.UvicornWorker -w 4 --bind 0.0.0.0:8000 --timeout 120 --log-level debug
